# app/supervisord.conf
[supervisord]
nodaemon=true

[program:nodejs]
command=pm2-runtime start /app/download.js -i max
autostart=true
autorestart=true
stderr_logfile=/var/log/nodejs.err.log
stdout_logfile=/var/log/nodejs.out.log

[program:php]
command=/bin/bash /app/startWeb.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/php.err.log
stdout_logfile=/var/log/php.out.log